# deploy/com.renquant.stops-liveness.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Software-stop registry liveness pager — S-FRAC stage-3 ops (#471 operator
  shortlist item 2: "pager scheduled nowhere, page never test-fired").

  Runs scripts/stops_liveness_pager.sh every 10 minutes. The wrapper invokes
  the EXECUTION-repo liveness checker (renquant_execution.software_stops_liveness,
  renquant-execution#29 — moved out of the deprecated umbrella per the
  Codex review of this package's prior revision), resolving the pinned
  execution/pipeline/common checkouts through the R-PIN Stage-1 runtime
  inventory (~/.renquant/deploy/runtime-inventory.json, read via
  renquant_orchestrator.deployment_manifest.load_runtime_inventory — the
  neutral per-host repo-name to checkout-path map, no umbrella lock-file
  dependency), and pages the LIVE ops ntfy topic on STALE / CORRUPT /
  checker-crash / inventory-resolution failure. The checker
  itself returns OK off-session, so a simple all-day StartInterval is the
  robust schedule: market-hours gating lives in ONE place (the checker),
  not duplicated here.

  TEMPLATE ONLY: merging this file does NOT install anything. Arming is a
  separately-granted operator landing step (the "apply" flag is two hyphens;
  XML comments cannot contain a double hyphen):

      scripts/install_stops_pager.sh install [apply flag]
      scripts/install_stops_pager.sh test-fire STALE   (delivery + response drill)

  BLOCKED — NOT MERGEABLE AS "DONE" (Codex round-3 review, 2026-07-11): this
  package is CORRECT for what it covers (pinned schedule + notification
  consumer, code resolution through the R-PIN runtime inventory) but it is
  NOT COMPLETE and must stay in changes-requested/blocked status. Two
  independent, load-bearing gaps, neither closed by this file:

  (1) DATA-ROOT AUTHORITY (the still-open Codex blocker): the
      RENQUANT_STOPS_PAGER_DATA_ROOT value below still names the deprecated
      umbrella checkout — a real production dependency, not a harmless
      default, even though code imports resolve through pins. This repo
      defines the neutral contract the writer should eventually satisfy
      (renquant_orchestrator.software_stops_registry_contract — neutral
      runtime-state-root convention + a versioned registry-file envelope)
      and the wrapper now WARNS on every run when the resolved data root is
      not that neutral root (it is, today). Migrating the ACTUAL writer
      (the live sell-only loop) off the umbrella path is a SEPARATE,
      explicitly-authorized, live-tree R-PIN landing change — not made
      here, not implied by merging this file.
  (2) SLA HONESTY: the alert-latency envelope is ~18-28 minutes after the
      first missed pass, which does NOT meet the design's 15-minute
      target.

  See doc/progress/2026-07-11-stops-liveness-pager-package.md, "BLOCKING
  FOLLOW-UP", for the full gap analysis and both required follow-ups. This
  package does NOT by itself support strategy-104#55/#56 enablement.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.renquant.stops-liveness</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/renhao/git/github/renquant-orchestrator/scripts/stops_liveness_pager.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>600</integer>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/renhao/.renquant/ops/stops-liveness/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/renhao/.renquant/ops/stops-liveness/launchd.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <!-- Explicit operator configuration (deploy-plist convention set by
             Codex r2 on #460 / deploy/com.renquant.shadow-ab-daily.plist):
             the paging topic is configuration owned by this plist and
             reviewed at arming time. "renquant" is the LIVE ops topic — the
             same channel the live sell-only loop (umbrella
             scripts/intraday_sell_104.sh NTFY_TOPIC) alerts on. -->
        <key>RENQUANT_STOPS_PAGER_NTFY_TOPIC</key>
        <string>renquant</string>
        <!-- RUNTIME CONTRACT: the wrapper script has NO default interpreter
             or data root (Codex r2 on #460 / this package's 2026-07-11
             re-review) — both are explicit, reviewed arming-time
             configuration here, never a hardcoded default in the script. -->
        <!-- Same interpreter the live pipeline stack already runs under
             (the conda env RenQuant/.venv/bin/python symlinks to),
             referenced directly, NOT through the umbrella .venv path.
             Verified live 2026-07-11: the BARE conda interpreter suffices
             because the wrapper puts the checker's full first-party import
             closure (common/base-data/artifacts/model/pipeline/execution)
             on PYTHONPATH from the runtime inventory - only the third-party
             scientific stack comes from this env; none of the umbrella
             venv's editable installs are needed. -->
        <key>RENQUANT_STOPS_PAGER_PYTHON</key>
        <string>/Users/renhao/miniconda3/envs/renquant/bin/python</string>
        <!-- INTERIM / BRIDGE VALUE — STILL THE DEPRECATED UMBRELLA (Codex
             round-3 review, 2026-07-11: "show the deployed plist/wrapper
             value and prove it does not name the umbrella" — it DOES; this
             comment does not dispute that). Where the live sell-only loop
             currently writes the software-stop registry. This value is
             explicit, reviewed plist configuration, never a wrapper-script
             default (RUNTIME CONTRACT) — but it is still a real production
             dependency on the umbrella. Migrating the WRITER to the neutral
             runtime-state-root contract this repo now defines
             (renquant_orchestrator.software_stops_registry_contract) is a
             SEPARATE, out-of-scope, live-tree R-PIN landing change; the
             wrapper detects and WARNS on this exact interim state every
             run (see scripts/stops_liveness_pager.sh, "STILL-OPEN
             BLOCKER"). -->
        <key>RENQUANT_STOPS_PAGER_DATA_ROOT</key>
        <string>/Users/renhao/git/github/RenQuant</string>
    </dict>
</dict>
</plist>
