# deploy/com.renquant.crypto-session.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  D-C11 crypto session scheduler — paper-mode tick every 15 minutes (24/7).

  TEMPLATE ONLY: merging this file does NOT install anything. Arming is a
  separately-granted operator landing step:

      1. Set RENQUANT_CRYPTO_TRADING=1 in the umbrella .env
      2. Optionally place a config JSON at the path below
      3. mkdir -p ~/Library/Logs/renquant
         cp deploy/com.renquant.crypto-session.plist ~/Library/LaunchAgents/
         launchctl load ~/Library/LaunchAgents/com.renquant.crypto-session.plist

  NOTE: entries are currently structurally blocked by the trust anchor gate
  (ENTRY_AUTHORIZATION_TRUST_ANCHOR_READY=False in crypto_session.py).
  The scheduler will run, evaluate gates, and log results, but will not
  authorize any entries until the trust anchor is wired.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.renquant.crypto-session</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/renhao/git/github/RenQuant/.venv/bin/python</string>
        <string>/Users/renhao/git/github/renquant-orchestrator/scripts/crypto_session_runner.py</string>
        <string>--once</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/renhao/git/github/RenQuant</string>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>StandardOutPath</key>
    <string>/Users/renhao/Library/Logs/renquant/crypto-session.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/renhao/Library/Logs/renquant/crypto-session.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/renhao/git/github/renquant-orchestrator/src:/Users/renhao/git/github/renquant-execution/src:/Users/renhao/git/github/renquant-common/src</string>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
