# ops/renquant105/com.renquant.rq105-session-scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.renquant.rq105-session-scheduler</string>
  <key>ProgramArguments</key><array><string>/bin/zsh</string><string>/Users/renhao/git/github/renquant-orchestrator-run/ops/renquant105/run_session_scheduler.sh</string></array>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>25</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>25</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>25</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>25</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>25</integer></dict>
  </array>
  <key>StandardOutPath</key><string>/Users/renhao/git/github/RenQuant/logs/rq105/launchd_session-scheduler.out</string>
  <key>StandardErrorPath</key><string>/Users/renhao/git/github/RenQuant/logs/rq105/launchd_session-scheduler.err</string>
</dict>
</plist>
